# .github/scripts/java_staging_app.service
[Unit]
Description=Java Boilerplate Staging
After=network.target

[Service]
User=teamjava
Group=teamjava
WorkingDirectory=/home/teamjava/boilerplate/staging/target
ExecStart=/usr/bin/java -jar /home/teamjava/boilerplate/staging/target/hng-java-boilerplate-0.0.1-SNAPSHOT.jar

[Install]
WantedBy=multi-user.target
